# systemd/pi-console.service
[Unit]
Description=pi-console API
After=network.target

[Service]
User=kentaro
WorkingDirectory=/home/kentaro/work/pi-console
EnvironmentFile=/home/kentaro/work/pi-console/.env
ExecStart=/usr/bin/python3 -m api.main
Restart=always

[Install]
WantedBy=multi-user.target
